# extras/supervisor/pyrobotdhs.ini
[program:robot-dhs]
user=blctl
command=/srv/pyrobotdhs/.venv/bin/pyrobotdhs --config /srv/pyrobotdhs/config.json --dcss 10.109.3.121
directory=/srv/pyrobotdhs
environment=PYTHONUNBUFFERED=1,BEAMLINE=Cage1
stopasgroup=true
